# Units/demo.service
[Unit]
Description=Demo Dienst - demonstriert einen ersten Dienst

[Service]
Type=oneshot
ExecStart=echo "Hello form Demo Service"