# Units/heartbeat-ef.socket
[Unit]
Description=Heartbeat Erfurt socket Activation by socket

[Socket]
ListenStream=666
#Service=heartbeat_ef.service