# Units/heartbeat-ef.service
[Unit]
Description=Show Erfurt heartbeat

[Service]
Type=oneshot
ExecStart=/usr/local/bin/heartbeat erfurt
